# hasiimusic.service
[Unit]
Description=Hasii Music Bot - Telegram Music Streaming Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/HasiiMusicBot
Environment="PATH=/root/HasiiMusicBot/hasiimusic/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/HasiiMusicBot/hasiimusic/bin/python3 -m HasiiMusic
Restart=always
RestartSec=10
StandardOutput=append:/root/HasiiMusicBot/log.txt
StandardError=append:/root/HasiiMusicBot/log.txt

[Install]
WantedBy=multi-user.target
